# hologram.service
[Unit]
Description=Hologram

[Service]
Type=simple
ExecStart=/home/cam/.cargo/bin/hologram start
User=cam

[Install]
WantedBy=multi-user.target
